# packages/genocide/genocide-26.tar.gz/genocide-26/genocide.service
[Unit]
Description=GENOCIDE - the king commits genocide
After=network-online.target

[Service]
DynamicUser=True
StateDirectory=genocide
LogsDirectory=genocide
CacheDirectory=genocide
ExecStart=/usr/local/bin/genocided
CapabilityBoundingSet=CAP_NET_RAW

[Install]
WantedBy=multi-user.target
