# packages/isc-py-common/isc-py-common-0.0.300.tar.gz/isc-py-common-0.0.300/isc_common/deploy/daphne-dev-SITENAME-channels.service
[Unit]
Description=Daphne server for SITENAME

[Service]
Restart=on-failure
User=USERNAME
Environment=EMAIL_PASSWORD=SECRET SMTP_EMAL_PASSWORD=SMTP_SECRET
WorkingDirectory=/home/USERNAME/sites/SITENAME/source
ExecStart=/home/USERNAME/.local/bin/daphne -b 0 --access-log ../accsee.log -p 8000 PROJECT_DIR_NAME.asgi-dev:application

[Install]
WantedBy=multi-user.target
