# packages/isc-py-common/isc-py-common-0.0.300.tar.gz/isc-py-common-0.0.300/isc_common/deploy/gunicorn-SITENAME.service
[Unit]
Description=Gunicorn server for SITENAME

[Service]
Restart=on-failure
User=USERNAME
Environment=EMAIL_PASSWORD=SECRET SMTP_EMAL_PASSWORD=SMTP_SECRET
WorkingDirectory=/home/USERNAME/sites/SITENAME/source
ExecStart=/home/USERNAME/.local/bin/gunicorn --bind unix:/tmp/SITENAME.socket --workers WORKERS --capture-output --access-logfile ../access.log --error-logfile ../error.log PROJECT_DIR_NAME.wsgi:application --timeout 30

[Install]
WantedBy=multi-user.target
