# packages/amz-tool/amz_tool-2021.2.1-py3-none-any.whl/amz/data/auto_update/amz_updater.service
[Unit]
Description=AMZ Auto Updater Service

[Service]
User=__USER__
Type=oneshot
ExecStart=/home/__USER__/.amz/data/auto_update/amz_update.py

[Install]
WantedBy=multi-user.target
